# backend/scripts/supervisor.bt.example.ini
; 宝塔 → 软件商店 → Supervisor 管理器 → 添加守护进程
; 将路径改为服务器上的实际目录，例如 /www/wwwroot/toursimagent/backend

[program:xhs-sign]
command=/www/server/pyporject_evn/versions/3.11.9/bin/python3 scripts/xhs_sign_server.py
directory=/www/wwwroot/toursimagent/backend
autostart=true
autorestart=true
startsecs=3
startretries=5
user=www
environment=XHS_SIGN_PORT="8765"
stdout_logfile=/www/wwwlogs/xhs-sign.log
stderr_logfile=/www/wwwlogs/xhs-sign.err.log
stopasgroup=true
killasgroup=true
